# quagga/supervisord.conf
[supervisord]
nodaemon=true

[program:zebra]
command=/usr/lib/quagga/zebra
[program:bgpd]
command=/usr/lib/quagga/bgpd
