# quagga/supervisord.conf
[supervisord]
nodaemon=true

[program:zebra]
command=/usr/sbin/zebra
[program:bgpd]
command=/usr/sbin/bgpd
